# packaging/rpm/content/script/prerm
#!/bin/sh

# Exit on non-zero exit code
set -e

# Unregister service
chkconfig --del motech

# Remove fragment dir if empty
FRAGMENT_DIR=/usr/share/motech/.motech/bundles/fragments

if [ -d $FRAGMENT_DIR ] && [ ! "$(ls -A $FRAGMENT_DIR)" ]; then
    rm -rf $FRAGMENT_DIR
fi

# Delete the motech user, if he exists
if [ ! `grep -c motech: /etc/passwd` -eq 0 ]; then
    skill -KILL -u motech
    userdel motech
fi

#Remove symlinks
rm -f /usr/share/motech/logs
rm -f /usr/share/motech/temp
rm -f /usr/share/motech/work
rm -f /usr/share/motech/webapps
rm -f /usr/share/motech/.motech
rm -f /usr/share/motech/felix-cache
